# deploy/supervisor/quantitative-trading-scheduler.conf
[program:quantitative-trading-scheduler]
directory=/root/quantitative_trading/backend
command=/root/quantitative_trading/deploy/start_scheduler.sh
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/quantitative_trading/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PYTHONUNBUFFERED="1"

